# packaging/systemd/nomwatch.service
[Unit]
Description=NomWatch local-first host
After=network-online.target
Wants=network-online.target nomwatch-tailscale-helper.service

[Service]
Type=simple
User=nomwatch
Group=nomwatch
ExecStart=/usr/bin/nomwatch host
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
KillMode=control-group
UMask=0077
StateDirectory=nomwatch
RuntimeDirectory=nomwatch
Environment=NOMWATCH_HOME=/var/lib/nomwatch
Environment=NOMWATCH_RUNTIME_DIR=/run/nomwatch
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nomwatch /run/nomwatch
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
